# scheduler.service
# SETUP INSTRUCTIONS:
# 1. Replace YOUR_USERNAME with your actual username
# 2. Replace /path/to/scheduler.py with the actual path to your scheduler.py file
# 3. Replace /path/to/scheduler/directory with the directory containing scheduler.py
#
# To install the service:
# 1. Copy this file to /etc/systemd/system/scheduler.service
# 2. sudo systemctl daemon-reload
# 3. sudo systemctl enable scheduler
# 4. sudo systemctl start scheduler
#
# VERIFICATION AND MONITORING:
#
# 1. Check Service Status:
# sudo systemctl status scheduler
#   - Should show "active (running)"
#   - Check for any error messages
#   - Verify the correct user is running the service
#
# 2. View Real-time Logs:
# journalctl -u scheduler -f
#   - Shows live log output
#   - Look for successful command executions
#   - Monitor for any error messages
#
# 4. Verify Process Running:
# ps aux | grep scheduler.py
#   - Should show the Python process running
#   - Verify correct user and working directory
#
# 5. Check System Resources:
# top -p $(pgrep -f scheduler.py)
#   - Monitor CPU and memory usage
#   - Check for any resource issues
#
# 6. Test Service Recovery:
#
# To stop the service:
# sudo systemctl stop scheduler
#
# To disable automatic startup:
# sudo systemctl disable scheduler
#
# To restart the service:
# sudo systemctl restart scheduler

[Unit]
Description=MongoDB Scheduler Service

[Service]
Type=simple
User=kjaved
# Replace with the absolute path to your scheduler.py
ExecStart=/home/kjaved/pytorch_env/bin/python /home/kjaved/mongodb-scheduler/scheduler.py
WorkingDirectory=/home/kjaved/mongodb-scheduler/
# Automatically restart on failure
Restart=always
# Wait 10 seconds before restarting
RestartSec=10
# Maximum number of restarts within the time window
StartLimitInterval=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target

